# ops/com.myandong.42space-event-arm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myandong.42space-event-arm</string>

  <key>WorkingDirectory</key>
  <string>/Users/myandong/Projects/42space</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="/Users/myandong/.nvm/versions/node/v22.22.3/bin:$PATH"; exec env DRY_RUN=0 EXECUTE=1 WAIT_FOR_RECEIPT=1 EXECUTION_RETRY_MS=500 EVENT_OPEN_WINDOW_SECONDS=60 I_UNDERSTAND_42_PRICE_MARKET_RISK=YES I_AM_NOT_IN_RESTRICTED_JURISDICTION=YES STAKE_PER_OUTCOME_USDT=2 EVENT_OUTCOME_COUNT=5 ARM_WAIT_FOR_FUNDING=1 ARM_FUNDING_RETRY_MS=60000 ARM_FUNDING_HOT_RETRY_MS=1000 ARM_FUNDING_HOT_WINDOW_MS=600000 ARM_CATCH_UP_AFTER_FUNDING=1 ARM_CATCH_UP_WINDOW_MS=60000 npm run event:arm</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/myandong/Projects/42space/data/logs/launch-event-arm.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/myandong/Projects/42space/data/logs/launch-event-arm.err.log</string>
</dict>
</plist>
